# dist/nix-serve.service
[Unit]
Description=nix binary cache creation tool

[Service]
User=nix-serve
Environment=NIX_DAEMON=remote
ExecStart=/usr/bin/nix-serve -p 8080

